# lcdgps/lcdproc.service
[Unit]
Description=lcdproc client
Requires=lcdd.service
After=lcdd.service

[Service]
Type=forking
ExecStart=/usr/local/bin/lcdproc -c /usr/local/etc/lcdproc.conf

[Install]
WantedBy=multi-user.target

